# supervisor_conf/suspend-lid-closed.conf
[program:suspend-lid-closed]
command = /usr/local/sbin/suspend-lid-closed.py
user = root
